# clash/supervisor/clash.ini
[program:clash]
command=$EXE
loglevel=warn
stopwaitsecs=1
stdout_logfile=/var/log/supervisor/clash.log
stderr_logfile=/var/log/supervisor/clash.log
